# .config/systemd/user/notify-time.service
[Unit]
Description=Notify User of Current Time

[Service]
Type=oneshot
ExecStart=/home/jalio/.local/bin/notify-time
