# dot_config/systemd/user/notify-time.service
[Unit]
Description=Notify User of Current Time

[Service]
Type=oneshot
ExecStart=notify-time