# dot_config/systemd/user/xsession.target
[Unit]
Description=Xsession running
#BindsTo=default.target
